# services/squeezelite.service
[Unit]
Description=Squeezelite (Lyrion)
After=network.target sound.target

[Service]
Type=simple
User=lyrion
ExecStart=/usr/bin/squeezelite -n "Lyrion DAC" -o hw:1
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target